# linux/crontab
0 0 * * 1-5 echo "#$(date +'\%s') TODAY  $(date +'\%a \%b \%d \%Y')" >> ~/.bash_history
0 8,12,18 * * 1-5 echo "#$(date +'\%s') TIME   $(date +'\%a \%b \%d \%Y \%R:\%S \%Z')" >> ~/.bash_history
